# CortanaKernel/Scripts/cortana-web.service
[Unit]
Description=Cortana Web Dashboard
After=network.target cortana-kernel.service
Requires=cortana-kernel.service
PartOf=cortana-kernel.service

[Service]
Type=simple
WorkingDirectory=%h/Cortana
ExecStart=%h/Cortana/CortanaWeb/out/CortanaWeb
ExecStartPre=/usr/bin/dotnet build %h/Cortana/CortanaWeb -o %h/Cortana/CortanaWeb/out --artifacts-path %h/Cortana/CortanaWeb/out/lib
EnvironmentFile=%h/.config/Cortana/Config/.env
Restart=on-failure
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=CortanaWeb
