# CortanaKernel/Scripts/cortana-telegram.service
[Unit]
Description=Cortana Telegram Bot
After=network.target redis.service cortana-kernel.service
Requires=cortana-kernel.service
PartOf=cortana-kernel.service

[Service]
Type=simple
WorkingDirectory=%h/Cortana
ExecStart=%h/Cortana/CortanaTelegram/out/CortanaTelegram
ExecStartPre=/usr/bin/dotnet build %h/Cortana/CortanaTelegram -o %h/Cortana/CortanaTelegram/out --artifacts-path %h/Cortana/CortanaTelegram/out/lib
EnvironmentFile=%h/.config/cortana/.env
Restart=on-failure
RestartSec=5
KillSignal=SIGTERM

StandardOutput=journal
StandardError=journal
SyslogIdentifier=CortanaTelegram
