# CortanaKernel/Scripts/cortana-kernel.service
[Unit]
Description=Cortana Kernel
After=network.target redis.service
Wants=redis.service cortana-discord.service cortana-telegram.service cortana-web.service

[Service]
Type=simple
WorkingDirectory=%h/Cortana
ExecStart=%h/Cortana/CortanaKernel/out/CortanaKernel
ExecStartPre=/usr/bin/dotnet build %h/Cortana/CortanaKernel -o %h/Cortana/CortanaKernel/out --artifacts-path %h/Cortana/CortanaKernel/out/lib
EnvironmentFile=%h/.config/cortana/.env
Restart=on-failure
RestartSec=5
KillSignal=SIGTERM

StandardOutput=journal
StandardError=journal
SyslogIdentifier=CortanaKernel

[Install]
WantedBy=default.target
